# ops/systemd/sost-alpha-status-export.service
[Unit]
Description=SOST Alpha Status Export
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/sost/SOST/sost-gold-exchange-private
ExecStart=/usr/bin/python3 scripts/export_alpha_status_json.py --output /opt/sost/website/api/alpha_live_status.json
