# ops/systemd/sost-alpha-positions-export.timer
[Unit]
Description=Export SOST positions every 5 minutes

[Timer]
OnBootSec=30
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
