# ops/systemd/sost-alpha-deals-export.service
[Unit]
Description=SOST Alpha Deal Export
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/sost/SOST/sost-gold-exchange-private
ExecStart=/usr/bin/python3 scripts/export_deals_live_json.py --output /opt/sost/website/api/deals_live.json
